# .config/autostart/gnome-startup.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=GNOME Startup
Exec=/home/leo/dotfiles/scripts/gnome-startup.sh
